# autostart/chromium_autostart.desktop
[Desktop Entry]
Type=Application
Name=LXRandR autostart
Comment= Chromium autostart
Exec=/home/tvpc/kiosk.sh
OnlyShowIn=LXDE
